# server-config/systemd/mockmail-email-processor.service
# DEPRECATED: This systemd service is no longer used in production.
# The email processor now runs as a TypeScript process managed by PM2:
#   pm2 start ecosystem.config.js
#   Process name: mockmail-processor
#
# This file is kept for reference only.

[Unit]
Description=MockMail Email Processor (DEPRECATED - use PM2)
Documentation=https://github.com/DiegoNogueiraDev/mockmail.dev
After=network-online.target postfix.service docker.service mockmail-api.service
Wants=network-online.target
Requires=postfix.service docker.service

[Service]
Type=simple
User=mockmail
Group=mockmail
WorkingDirectory=/home/mockmail/app/backend
ExecStart=/usr/bin/node /home/mockmail/app/backend/dist/emailProcessor.js
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-processor
TimeoutStartSec=60

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/mockmail
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
